# socksproxy.service
[Unit]
Description=Proxy between the enclave and host
After=network.target

[Service]
ExecStart=/opt/nitro-shim/socksproxy
Environment=SOCKS_PROXY_LISTEN_ADDR=127.0.0.1:1080
Environment=SOCKS_PROXY_ALLOWED_ADDRS=
Environment=SOCKS_PROXY_ALLOWED_FQDNS=acme-v02.api.letsencrypt.org,nitro.nymity.ch

[Install]
WantedBy=multi-user.target
